# deploy/app/crontab
# supercronic crontab — times are in the container TZ (America/New_York, set in compose) so
# the schedule is US-market-relative regardless of DST. Fields: min hour dom mon dow.
#
# CADENCE (decided): weekly scan + quarterly rebalance. The Tedium Premium is a patient,
# filing-driven strategy, so the weekly job SCANS the week's new EDGAR filers and refreshes
# the candidate book + alerts; the actual book is rebalanced on a slower cadence (see below).
#
# LIVE (paper): the weekly job TRANSMITS orders to the IBKR paper account. It runs the L3 MD&A
# Deterioration Lead under a HARD weekly cap (--max-llm-usd 2) and rebalances toward the BUY book
# (long-only, policy-weighted; BUY needs L3 clearance via the safety gate). Monday 10:00 ET = 30
# min after the open so market orders fill. Revert to preview = drop --transmit.
0 10 * * 1 cd /app && python -m deepvalue.forward.run --execute ibkr --transmit --max-llm-usd 2 --deepdive-usd 15 >> /app/data/forward/logs/cron.log 2>&1

# Daily watchdog 08:00 ET: alerts (Telegram/email) if the last run is missing/errored/stale
# (>8 days), so a silently-skipped Monday is caught within a day.
0 8 * * * cd /app && python -m deepvalue.forward.run --healthcheck >> /app/data/forward/logs/healthcheck.log 2>&1
